# CP3/config-nginx.service
[Unit]
Description=Configure Nginx server

[Service]
Type=oneshot
ExecStart=/home/ubuntu/nginx-cloud-init.sh

[Install]
WantedBy=multi-user.target
